# build/usr/lib/opkg/info/MtpDaemon.list
/etc/init.d/mtp
/usr/sbin/MtpDaemon
/usr/bin/MtpTools
